# misc/zulip-mail-mirror.timer
[Unit]
Description=Mirrors new emails to Zulip every minute.

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target